# files/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
user=nobody
command=/usr/sbin/haproxy -f %(ENV_CONF_HAPROXY)s
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
